# templates/uzblmonitor.service.erb
[Unit]
Description=uzblmonitor service
After=nodm-uzblmonitor.service

[Service]
Environment="DISPLAY=:0"
User=monitor
ExecStart=/usr/bin/uzblmonitor <%= @browser %>
Restart=always

[Install]
WantedBy=graphical.target
